# yure-tweet-checker/supervisord.conf
[program:node]
command = node yure-tweet-checker.js
directory = /home/dotcloud/current
