# yure-tweet-express/supervisord.conf
[program:node]
command = node yure-tweet-express.js
directory = /home/dotcloud/current
