# yure-tweet-crawler/supervisord.conf
[program:node]
command = node yure-tweet-crawler.js
directory = /home/dotcloud/current
